# root/etc/services.d/socat/run
#!/usr/bin/with-contenv sh

umask 002

s6-echo "Forwarding $LISTEN_PORT -> $FORWARD_HOST:$FORWARD_PORT"

exec socat tcp-listen:${LISTEN_PORT},reuseaddr,fork tcp:${FORWARD_HOST}:${FORWARD_PORT}
